# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/hotplug.d/usb/70-time_machine
#!/bin/sh

# echo ACTION=$ACTION, DEVPATH=$DEVPATH, SUBSYSTEM=$SUBSYSTEM, BUSNUM=$BUSNUM >/dev/console

case "$ACTION" in
        add)
            ps |grep "sh /lib/time_machine/tm_usb_add.sh" |grep -v grep |awk '{print $1}' |xargs kill -9
            sh /lib/time_machine/tm_usb_add.sh &
            ;;
        remove)
            sh /lib/time_machine/tm_usb_remove.sh
            ;;
esac
