# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/hotplug.d/cable/00-autodial
#!/bin/sh

device=`uci get profile.@wan[0].wan_ifname -c /etc/profile.d/`
[ "$device" = "$DEVICE" -a up = "$ACTION" ] && {                                                                
    local autodial_flag=`uci get network.wan.autodial_flag`
    [ $autodial_flag -eq 1 ] && [ ! -e /tmp/reboot ] && {
        echo autodial begin !!!! > /dev/console                                     
        wan_connected && lua /lib/autodial/auto_dial.lua auto_dial &
    }                                                           
}
