# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/hotplug.d/firewall/22-yandex_dns
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

if [ "$ACTION" = "start" ]; then
	local sysmode
	config_load sysmode
	config_get sysmode sysmode mode "router"
	local country=`getfirm SPECIAL_ID`
	if [ "$country" != "52550000" -o "$sysmode" = "ap" ]; then
		return 0
	fi
	/etc/init.d/yandex_dns restart
fi
