# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/init.d/imb
#!/bin/sh /etc/rc.common

START=66
#STOP=66

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	local sysmode
	config_load sysmode
	config_get sysmode sysmode mode "router"
	if [ "$sysmode" = "router" ]; then
	service_start /usr/sbin/imbd
	fi

	return 0
}

stop() {
	local sysmode
	config_load sysmode
	config_get sysmode sysmode mode "router"
	if [ "$sysmode" = "router" ]; then
	service_stop /usr/sbin/imbd
	fi

	return 0
}
