# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/init.d/luarsa_keys_gen
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=48

start() {
    #if [ ! -e usr/bin/openssl ]; then
    #    luarsa_keys_gen
    #fi
}

restart() {
    #if [ ! -e usr/bin/openssl ]; then
    #   luarsa_keys_gen
    #fi
}

reload() {
    #if [ ! -e usr/bin/openssl ]; then  
    #   luarsa_keys_gen
    #fi
}
