# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/init.d/network
#!/bin/sh /etc/rc.common

START=25
STOP=90

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

NETWORK_MOD_ID=290
LOGX_ID_NETIFD_START=200
LOGX_ID_NETIFD_STOP=201

start() {
	local pid=""
	local enable
	config_load tfstats
	config_get enable "switch" enable

	# BCM4908: work around a flow cache bug
	if [ -d /proc/fcache/ ] ; then
		local internet_proto=`uci get network.internet.proto 2>/dev/null`
		local pptp_vpn=`uci get pptpd.pptpd.enabled 2>/dev/null`
		if [[ "$pptp_vpn" == "on" ]] && [[ "$internet_proto" == "pptp" || "$internet_proto" == "l2tp" ]] ; then
			fcctl config --gre 0
		else
			fcctl config --gre 1
		fi
	fi

	# Intel:config ppa
	if [ -d /proc/ppa/ ] ; then
		config_ppa
	fi	
	
	stop

	[ -e /proc/sys/kernel/core_pattern ] && {
		ulimit -c unlimited
		echo '/tmp/%e.%p.%s.%t.core' > /proc/sys/kernel/core_pattern
	}

	# netifd depends on "statistics" kern module
	# we insmod statistics in boot, so no need to do it here
	# insmod statistics
	if [ "$enable" = "on" ] ; then
		ubus call tfstats on || true
	fi
	service_start /sbin/netifd

	pid=$(pidof netifd)
	if [ x"$pid" != x ]; then
		logx -p $pid $NETWORK_MOD_ID $LOGX_ID_NETIFD_START
	fi
}

restart() {
	start
}

shutdown() {
	stop
}

stop() {
	local pid=$(pidof netifd)

	ifdown -a
	if [ x"$pid" != x ]; then
		sleep 1
	fi

	service_stop /sbin/netifd
	ubus call tfstats off || true
	#rmmod statistics

	if [ x"$pid" != x ]; then
		logx -p $pid $NETWORK_MOD_ID $LOGX_ID_NETIFD_STOP
	fi
}

reload() {
	ubus call network reload
}
