# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/init.d/sync-server
#!/bin/sh /etc/rc.common

START=51

SERVICE_DAEMONIZE=1
start()
{
	local sysmode
	config_load sysmode
	config_get sysmode sysmode mode "router"

	local onemesh_enable=$(uci get onemesh.onemesh.enable)
	if [ -z "$onemesh_enable" ]; then
		onemesh_enable="on"
	fi

	if [ "$sysmode" = "router" -a "$onemesh_enable" = "on" ]; then
		service_start "/bin/nice" -n -5 /usr/bin/sync-server
	fi
}

stop()
{
	service_stop /usr/bin/sync-server
}
