# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/init.d/time_machine
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2012 OpenWrt.org

START=99

start()
{
	# wait for usb storage mounted with enough time 60s 
	# to check whether start time machine

	flock -x -w 10 "/var/run/tm_checkrun.lock" -c "sleep 60;/usr/sbin/tm_checkrun" &
}

