# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/init.d/wireless_schedule
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=48

WIRELESS_SCHEDULE_LIBDIR=/lib/wifi

wireless_schedule() {
	. ${WIRELESS_SCHEDULE_LIBDIR}/wireless_schedule_core.sh
	wireless_schedule_$1
}

start() {
	wireless_schedule start
}

stop() {
	wireless_schedule stop
}

restart() {
	wireless_schedule restart
}
