# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/init.d/zzzapmode
#!/bin/sh /etc/rc.common

START=99

start () {
    config_load sysmode
    config_get mode "sysmode" "mode"
    config_load network
    config_get ipaddr "wan" "ipaddr"

    if [ $mode == 'ap' ]; then
		/etc/init.d/smartip restart
    fi
    return 0
}
