# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/hotplug.d/iface/91-sysntpd
#!/bin/sh

if [ "$ACTION" != "ifup" ];then
	exit 0
fi

if [ "$INTERFACE" = "wan" -o "$INTERFACE" = "internet"  -o "$INTERFACE" = "mobile" ];then
	. /lib/functions.sh
	SYS_SET_TIME=$(uci_get_state systime core sync)
	if [ "$SYS_SET_TIME" -eq 0 ];then
		pgrep ntpd && /etc/init.d/sysntpd restart >/dev/null
	fi
fi
