# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/hotplug.d/iface/98-cwmp_hotplug
#!/bin/sh

[ "ifup" = "$ACTION" ] && {
    if [ "wan" = "$INTERFACE" -o "internet" = "$INTERFACE" ]; then
        if [ $(ps w|grep -v grep|grep -w cwmp -c) -gt 0 ]; then
            /etc/init.d/cwmp reload generate_conf
            cwmp_cli -u wan
        else
            /etc/init.d/cwmp reload
        fi
    fi
}
