# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/hotplug.d/iface/98-cloud-brd
[ "$ACTION" == "ifup" ] && cloud_reconnect
