# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/rc.d/S90nf-shortcut
#!/bin/sh /etc/rc.common
#insmod nf-shortcut module

START=90

start() {
	insmod nf_sc
	return 0
}

restart() {
	return 0
}

stop() {
	return 0
}

reload() {
	return 0
}
