# Archer C6(US)_V3.20_220923/a6c6_v3-up-ver1-0-15-P1[20220923-rel62166]_2022-09-23_17.43.32/etc/rc.d/S99zzzzboot_done
#!/bin/sh /etc/rc.common

START=99zzzz
STOP=90


boot() {
	## set phy_auto_up flag to make phy up automatically when restart network
	echo 1 > /proc/mt7621/phy_auto_up
	touch /tmp/dut_bootdone
	echo 3 > /proc/sys/vm/drop_caches
	echo "drop caches" > /dev/console
}

start() {
    return 0
}

stop() {
    return 0
}
